# ablecloud
#!/bin/sh /etc/rc.common
# Copyright Abcloud Team
#@/etc/init.d/
START=99
start() {
	if [ -f "/tmp/AbcloudFirmware" ]; then
		 mv /tmp/AbcloudFirmware /usr/bin
	fi
	 /usr/bin/AbcloudFirmware
}

